# s6/protonvpn/finish
#!/bin/sh

# stop container
s6-svscanctl -t /etc/s6
